# services/am2302-sensor.timer
[Unit]
Description=AM2302 Temp and humidity Sensor timer

[Timer]
OnUnitActiveSec=5m
OnStartupSec=15s

[Install]
WantedBy=timers.target
